# supervisord.conf
[supervisord]
nodaemon=true
[program:kong]
command=/usr/local/bin/kong start
autorestart=unexpected
autostart=true